# server/deploy/grilling-sleek.service
[Unit]
Description=grilling-sleek server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/grilling-sleek
WorkingDirectory=/srv/grilling-sleek
Restart=on-failure
RestartSec=5

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=35

# Environment (OTel collector, log level, etc.)
Environment=RUST_LOG=info
# Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://localhost:4317

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/srv/grilling-sleek/data /srv/grilling-sleek/log

[Install]
WantedBy=multi-user.target
