# docker/s6/s6-rc.d/backend/run
#!/bin/sh
cd /app
exec ./grilling-sleek
