# docker/s6/s6-rc.d/caddy/run
#!/bin/sh
exec caddy run --config /app/Caddyfile
